# src/startup/systemd/ssp.service
[Unit]
Description=ssp

[Service]
Type=simple
ExecStart=/usr/bin/ssp.sh

[Install]
WantedBy=multi-user.target
